# com.forexgod.morningscan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.forexgod.morningscan</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/forexgod/Desktop/Glitch in Matrix/trading-ai-agent apollo/run_morning_scan.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/forexgod/Desktop/Glitch in Matrix/trading-ai-agent apollo/logs/morning_scan_launchd.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/forexgod/Desktop/Glitch in Matrix/trading-ai-agent apollo/logs/morning_scan_launchd_error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>WorkingDirectory</key>
    <string>/Users/forexgod/Desktop/Glitch in Matrix/trading-ai-agent apollo</string>
</dict>
</plist>
