# com.forexgod.glitch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.forexgod.glitch</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/forexgod/Desktop/Glitch in Matrix/trading-ai-agent apollo/start_system.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/forexgod/Desktop/Glitch in Matrix/trading-ai-agent apollo/logs/launchd_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/forexgod/Desktop/Glitch in Matrix/trading-ai-agent apollo/logs/launchd_stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/forexgod/Desktop/Glitch in Matrix/trading-ai-agent apollo</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
